# voice-chat.service
[Unit]
Description=OpenClaw Voice Chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hodaddy
WorkingDirectory=/home/hodaddy/voice-chat
ExecStart=/home/hodaddy/voice-chat/.venv/bin/python /home/hodaddy/voice-chat/server.py
Restart=on-failure
RestartSec=5
Environment=PATH=/home/hodaddy/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
